# Repository/deployment/provisioning/upstart/sentiment-analyzer.conf
description "Sentiment Analyzer"
author      "PS AIS"

console log

start on filesystem or runlevel [2345]
stop on shutdown

respawn
respawn limit 5 120
post-stop exec sleep 10

script
    export SENTIMENT_CONFIG="/etc/psais/sentiment-analyzer.yaml"
    exec /usr/bin/python3 /opt/psais/sentiment-analyzer/sentiment_analyzer.py
end script
